# config/vless-vpn.service
[Unit]
Description=VLESS+Reality VPN Management System
Documentation=https://github.com/your-repo/vless-vpn
After=network.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=forking
RemainAfterExit=yes
WorkingDirectory=/opt/vless
ExecStartPre=/usr/bin/docker-compose -f /opt/vless/docker-compose.yml pull -q
ExecStart=/usr/bin/docker-compose -f /opt/vless/docker-compose.yml up -d
ExecStop=/usr/bin/docker-compose -f /opt/vless/docker-compose.yml down
ExecReload=/usr/bin/docker-compose -f /opt/vless/docker-compose.yml restart xray

# Restart settings
Restart=on-failure
RestartSec=10
StartLimitBurst=3
StartLimitInterval=60

# Security settings
User=root
Group=root

# Environment
Environment=COMPOSE_PROJECT_NAME=vless-vpn
Environment=COMPOSE_HTTP_TIMEOUT=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vless-vpn

# Process management
TimeoutStartSec=300
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target