# config/vless-vpn.service
[Unit]
Description=VLESS+Reality VPN Telegram Bot Service
Documentation=https://github.com/XTLS/Xray-core
After=network.target docker.service
Wants=docker.service
PartOf=docker.service

[Service]
Type=simple
User=root
Group=root

# Working directory
WorkingDirectory=/opt/vless

# Environment
Environment=PYTHONPATH=/opt/vless/modules
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Environment file
EnvironmentFile=-/opt/vless/config/bot_config.env

# Main process
ExecStart=/usr/bin/python3 /opt/vless/modules/telegram_bot.py
ExecReload=/bin/kill -HUP $MAINPID

# Process management
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vless
ReadOnlyPaths=/opt/vless/config

# Resource limits
MemoryLimit=512M
TasksMax=100

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vless-vpn

# Kill settings
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target